# etc/monitor.service
[Unit]
Description=Monitor Kernel Module
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/modprobe hello
ExecStop=/sbin/rmmod hello
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target